# config/crontab
*/30 * * * * /usr/bin/python3 /home/supervision/supervision/analyse_and_alert.py --credentials /home/supervision/credentials.txt
